# 08-supervisor-intro/01-basic.conf
[program:01-basic]
command=php /srv/app/01-basic/consumer.php
