# 09-supervisor-multiple/01-basic.conf
[program:01-basic]
command=php /srv/app/01-basic/consumer.php
directory=/srv/app/
autorestart=true
process_name=%(program_name)s_%(process_num)02d
numprocs=3



