# 10-supervisor-autorestart/10-supervisor-autorestart.conf
[program:10-supervisor-autorestart]
command=php /srv/app/10-supervisor-autorestart/consumer.php
directory=/srv/app/
autorestart=true
process_name=%(program_name)s_%(process_num)02d
